# deploy/zimbra-import-worker.service
[Unit]
Description=Zimbra Import Worker
After=network.target

[Service]
User=zimbra-import
WorkingDirectory=/opt/zimbra-import
ExecStart=/opt/zimbra-import/venv/bin/python -m zimbra_import.worker /etc/zimbra-import/config.ini
Restart=on-failure

[Install]
WantedBy=multi-user.target
